# systemd/lyrebird-stream.service
# LyreBirdAudio Stream Manager Service
#
# This service manages audio streaming from USB devices to RTSP via MediaMTX.
# It automatically detects USB audio devices and maintains persistent streams
# with automatic restart on failure.
#
# Installation:
#   sudo cp lyrebird-stream.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable lyrebird-stream
#   sudo systemctl start lyrebird-stream
#
# Configuration:
#   Primary: /etc/lyrebird/config.yaml
#   Overrides: /etc/lyrebird/environment (optional)
#
# Logs: journalctl -u lyrebird-stream -f
#
# Hot-reload configuration (no restart required):
#   sudo systemctl reload lyrebird-stream
#
# User Configuration:
#   The service runs as root by default for reliable ALSA device access.
#   To run as a dedicated user (recommended for production):
#   1. Create user: sudo useradd -r -s /usr/sbin/nologin -G audio lyrebird
#   2. Create directories: sudo mkdir -p /var/run/lyrebird /etc/lyrebird
#   3. Set ownership: sudo chown lyrebird:audio /var/run/lyrebird
#   4. Change User=root to User=lyrebird in this file
#   5. Reload: sudo systemctl daemon-reload && sudo systemctl restart lyrebird-stream

[Unit]
Description=LyreBirdAudio Stream Manager
Documentation=https://github.com/tomtom215/lyrebirdaudio-go
After=network.target sound.target mediamtx.service
Wants=mediamtx.service
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=root
Group=audio

# Optional environment file for configuration overrides
# Create /etc/lyrebird/environment with:
#   LYREBIRD_CONFIG=/etc/lyrebird/config.yaml
#   LYREBIRD_LOG_LEVEL=info
EnvironmentFile=-/etc/lyrebird/environment

# Default configuration (can be overridden via environment file)
Environment=LYREBIRD_CONFIG=/etc/lyrebird/config.yaml
Environment=LYREBIRD_LOG_LEVEL=info

# Main executable
ExecStart=/usr/local/bin/lyrebird-stream --config=${LYREBIRD_CONFIG} --log-level=${LYREBIRD_LOG_LEVEL}

# Hot-reload configuration without stopping streams
ExecReload=/bin/kill -HUP $MAINPID

# Graceful shutdown
ExecStop=/bin/kill -SIGTERM $MAINPID
TimeoutStopSec=30

# Restart policy
Restart=always
RestartSec=10
# WatchdogSec is intentionally absent: the daemon does not call sd_notify(WATCHDOG=1),
# so enabling watchdog supervision would cause spurious service restarts (M-2).

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=yes
RestrictNamespaces=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
RestrictRealtime=yes
SystemCallFilter=@system-service
SystemCallArchitectures=native

# Allow access to audio devices
SupplementaryGroups=audio
DeviceAllow=/dev/snd/* rw
DevicePolicy=closed

# Allow access to required paths
ReadWritePaths=/var/run/lyrebird
ReadOnlyPaths=/etc/lyrebird /proc/asound

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# P-10 fix: Memory limits to prevent OOM on resource-constrained devices
# (e.g., Raspberry Pi with 1-4 GB RAM and multiple USB microphones).
# MemoryHigh triggers memory pressure reclaim; MemoryMax is a hard kill limit.
# 512M is conservative for audio-only FFmpeg streams (typically 20-50 MB each).
MemoryHigh=384M
MemoryMax=512M

# Environment
Environment=HOME=/var/run/lyrebird

[Install]
WantedBy=multi-user.target
